# supervisord.conf
[program:node]
command = node server/server.js
directory = /home/dotcloud/current